# clean-recently-used@.service.in
[Unit]
Description=Clean paths below %f from the recently used bookmarks

[Service]
Type=oneshot
ExecStart="$USER_BINARIES/clean-recently-used" %f
SystemCallFilter=@system-service
SystemCallArchitectures=native
